# systemd/picoclaw-gateway.service
[Unit]
Description=PicoClaw Gateway
After=network.target venice-proxy.service
Wants=venice-proxy.service

[Service]
Type=simple
User=pi
# The Telegram bot token is set via a systemd drop-in override, NOT in this file.
# To configure:
#   sudo systemctl edit picoclaw-gateway.service
# Then add:
#   [Service]
#   Environment="PICOCLAW_CHANNELS_TELEGRAM_TOKEN=your-telegram-bot-token"
# Or use an EnvironmentFile:
#   EnvironmentFile=/etc/picoclaw.env
ExecStart=/usr/local/bin/picoclaw gateway
Restart=always
RestartSec=5
WorkingDirectory=/home/pi

[Install]
WantedBy=multi-user.target