# systemd/venice-proxy.service
[Unit]
Description=Venice API Proxy
After=network.target

[Service]
Type=simple
User=pi
# Set your Venice API key here, or use an EnvironmentFile for better security:
# EnvironmentFile=/etc/venice-proxy.env
Environment="VENICE_API_KEY=your-venice-api-key-here"
ExecStart=/usr/bin/python3 /usr/local/bin/venice-proxy.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target