# systemd/venice-proxy.service
[Unit]
Description=Venice API Proxy
After=network.target

[Service]
Type=simple
User=pi
# The API key is set via a systemd drop-in override, NOT in this file.
# To configure:
#   sudo systemctl edit venice-proxy.service
# Then add:
#   [Service]
#   Environment="VENICE_API_KEY=sk-venice-xxxx"
# Or use an EnvironmentFile:
#   EnvironmentFile=/etc/venice-proxy.env
ExecStart=/usr/bin/python3 /usr/local/bin/venice-proxy.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target